# philips-tv.service
[Unit]
Description=Philips JointSpace Volume Queue Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/storage/tv-control
ExecStart=/usr/bin/python3 /storage/tv-control/philips_tv.py serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
